# supervisord.conf
[supervisord]
nodaemon=true

[program:couchpotato]
command=/usr/bin/python couch-potato/CouchPotato.py
